# ansible/roles/auto-mount-volumes/files/99-exosphere-automount-volume.rules
ACTION=="add|bind|change", KERNEL=="[sv]d[b-z]", ENV{SYSTEMD_WANTS}+="automount-volume@%k.service"
ACTION=="remove", KERNEL=="[sv]d[b-z]", RUN+="/usr/bin/systemctl stop automount-volume@%k.service"